# scripts/com.ivan.realestate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ivan.realestate</string>

    <key>WorkingDirectory</key>
    <string>/Users/ivancollazo/realestate-repo</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ivancollazo/realestate-repo/.venv/bin/realestate</string>
        <string>run</string>
        <string>--source</string>
        <string>all</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ivancollazo/realestate-repo/logs/realestate.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ivancollazo/realestate-repo/logs/realestate.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ivancollazo/realestate-repo/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>HOME</key>
        <string>/Users/ivancollazo</string>
    </dict>
</dict>
</plist>
